# private_dot_config/autostart/headsetcontrol-qt.desktop
[Desktop Entry]
Path=/usr/bin
Type=Application
Exec=/usr/bin/HeadsetControl-Qt
Name=HeadsetControl-Qt
