# private_dot_config/autostart/io.github.vikdevelop.SaveDesktop.server.desktop
[Desktop Entry]
Name=SaveDesktop (syncing server)
Type=Application
Exec=flatpak run io.github.vikdevelop.SaveDesktop --start-server