# private_dot_config/autostart/io.github.vikdevelop.SaveDesktop.sync.desktop
[Desktop Entry]
Name=SaveDesktop (syncing tool)
Type=Application
Exec=flatpak run io.github.vikdevelop.SaveDesktop --sync