# private_dot_config/autostart/moondeckbuddy.desktop
[Desktop Entry]
Type=Application
Name=MoonDeckBuddy
Exec=/var/home/zany130/AppImages/moondeckbuddy.appimage
Icon=MoonDeckBuddy
