# private_dot_config/autostart/moondeckbuddy.desktop
[Desktop Entry]
Type=Application
Name=MoonDeckBuddy
Exec=/var/home/zany130/AppImages/moondeckbuddy_a2f2a2.appimage
Icon=MoonDeckBuddy
