# private_dot_config/autostart/pcloud.desktop
[Desktop Entry]
Type=Application
Name=pCloud Drive
Comment=Cloud storage service for your files
Exec="/var/home/zany130/Applications/pcloud/pCloud.AppImage"
Icon=appimagekit-pcloud
StartupWMClass=pcloud
Terminal=false
Categories=Network;FileTransfer;
StartupNotify=false
X-GNOME-Autostart-enabled=true
