# private_dot_config/autostart/private_onedrivegui_5bb160.desktop
[Desktop Entry]
Categories=Utility;
Comment=
Exec=/home/zany130/AppImages/onedrivegui_5bb160.appimage
Icon=/home/zany130/AppImages/.icons/onedrivegui_5bb160.png
Name=OneDriveGUI
Terminal=false
TryExec=/home/zany130/AppImages/onedrivegui_5bb160.appimage
Type=Application
X-AppImage-Arch=x86_64
X-AppImage-Name=OneDriveGUI
X-AppImage-Version=1.0.3_fix116
