# private_dot_config/autostart/private_io.github.kotatogram.desktop
[Desktop Entry]
Exec=flatpak run --command=kotatogram-desktop io.github.kotatogram -autostart
Name=io.github.kotatogram
Type=Application
X-Flatpak=io.github.kotatogram
