# private_dot_config/systemd/user/plasma-plasmashell.service.d/override.conf
[Service]
ExecStartPre=/usr/bin/bash -c '/home/zany130/Scripts/killGamescopeApps.sh'
