# private_dot_config/systemd/user/alga.service
[Unit]
Description=turn on tv and set input

[Service]
Type=simple
ExecStart=/var/home/zany130/Scripts/wolTV.sh
