# private_dot_config/systemd/user/gamescopeApps.service
[Unit]
Description=headless apps to run in gamescope session

[Service]
Type=forking
ExecStart=/home/zany130/Scripts/startGamescopeAppsBazzite.sh
